# evolution_engine/service/evolution_engine.service
[Unit]
Description=Evolution Engine – Monitoring & Self-Evolution
After=network.target

[Service]
Type=oneshot
User=etherverse
ExecStart=/home/etherverse/etherverse/evolution_engine/scripts/run_engine.sh

[Install]
WantedBy=multi-user.target
